# arch/rootfs/usr/lib/systemd/system/nos-resource-monitor.service
[Unit]
Description=NeuralOS Resource Monitor Service
PartOf=nos.target
After=nos.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -c "from nos import NeuralNode; n = NeuralNode(); n.simulate_activity(5); print('Resource monitor ready')"
Restart=on-failure

[Install]
WantedBy=nos.target
