# arch/rootfs/usr/lib/systemd/system/nos.target
[Unit]
Description=NeuralOS Target
After=multi-user.target

[Install]
WantedBy=multi-user.target
